# net-p2p/jdownloader/files/jdownloader.initd

#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting jDownloader web interface"
	[ -x "/usr/bin/Xvfb" ] && /usr/bin/Xvfb :100 -ac >/dev/null 2>&1 &
	cd "${JD_HOME}"
	DISPLAY=":100.0" /usr/bin/java -Xmx768m -jar JDownloader.jar >/dev/null 2>&1 &
	eend $?
}

stop() {
	ebegin "Stopping jDownloader"
	killall Xvfb
	sleep 1
	eend $?
}
