# dev-db/apache-cassandra-bin/files/cassandra
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$
depend() {
    use net
}

PIDDIR=/var/run/pids
PIDFILE=$PIDDIR/cassandra.pid
[ -d $PIDDIR ] || mkdir -m=777 $PIDDIR


start() {
    ebegin "Starting Apache Cassandra"
    start-stop-daemon --start --quiet --user cassandra --pidfile $PIDFILE --exec /opt/cassandra-3.0.13/bin/cassandra -- -p $PIDFILE
    eend $?
}

stop() {
    ebegin "Stopping Apache Cassandra"
    start-stop-daemon --stop --quiet --pidfile $PIDFILE
    eend $?
}
