# net-misc/activemq-bin/files/init
#!/sbin/runscript

NAME=${SVCNAME/glassfish/}

depend() {
    need net
}

start() {
    ebegin "Starting Apache ActiveMQ application server"
    su -l activemq -c '/opt/activemq/bin/activemq start'
    eend $?
}

stop() {
    ebegin "Sopping GlassFish application server"
    su -l activemq -c '/opt/activemq/bin/activemq stop'
    eend $?
}
